# builder/overlay/etc/systemd/system/volt-display-hotplug.service
[Unit]
Description=Volt Display Manager — Hotplug Watcher
Documentation=https://samaris.dev/docs/volt-display-manager
After=nodm.service volt-display-manager.service
Wants=volt-display-manager.service

[Service]
Type=simple
User=user
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/user/.Xauthority
Environment=VOLT_DISPLAY_LASTGOOD_PATH=/home/user/.local/state/samaris/display.lastgood.toml
Environment=VOLT_DISPLAY_MANAGER_LOG=/tmp/volt-display-manager.log
ExecStart=/usr/local/bin/start-display-manager.sh --config-path /home/user/.local/state/samaris/display.generated.toml --event-path /home/user/.local/state/samaris/display.event.json --watch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
